# fence_virtd.service.in
[Unit]
Description=Fence-Virt system host daemon
After=syslog.target network.target libvirt-qmf.service libvirtd.service corosync.service
Requires=syslog.target network.target

[Install]
WantedBy=multi-user.target

[Service]
Type=forking
Restart=on-failure
SysVStartPriority=99
Environment="FENCE_VIRTD_ARGS=-w"
# Autogenerated below here
