# agents/virt/fence_virtd.service.in
[Unit]
Description=Fence-Virt system host daemon
Documentation=man:fence_virtd(8)
Documentation=man:fence_virt.conf(5)

After=basic.target
After=network.target
After=syslog.target
After=libvirtd.service
After=corosync.service

Requires=basic.target
Requires=network.target

[Install]
WantedBy=multi-user.target

[Service]
Type=forking
Restart=on-failure
Environment="FENCE_VIRTD_ARGS=-w"

# Autogenerated below here
